# roles/jg_moosefs_chunkserver/templates/moosefs-metalogger.service.j2
[Unit]
Description=MooseFS Metalogger server
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=simple
User=mfs
Group=mfs
ExecStart=/usr/sbin/mfsmetalogger start
ExecStop=/usr/sbin/mfsmetalogger stop
ExecReload=/usr/sbin/mfsmetalogger reload
PIDFile=/var/lib/mfs/.mfsmetalogger.lock
Restart=on-failure

[Install]
WantedBy=multi-user.target
